# .config/systemd/user/qotd.service
[Unit]
Description=Generate a quote of the day

[Service]
Type=oneshot
ExecStart=/home/humanwave/.local/bin/qotd

[Install]
WantedBy=default.target
