# .config/systemd/user/qotd.timer
[Unit]
Description=Generate a quote of the day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
