# .config/systemd/user/send-package.timer
[Unit]
Description=Do Not Ask

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
